# .config/supervisord.conf
[program:gunicorn]
command=gunicorn config.wsgi.production -c ../.config/gunicorn_cfg.py
directory=/srv/projects/app

[program:nginx]
command = nginx -g 'daemon off;'